# roles/ci_service/templates/etc/systemd/system/ci_service.service
[Unit]
Description=ci_service

[Service]
TimeoutStartSec=0
TimeoutStopSec=10
Restart=always
ExecStart=/opt/ci_service/ci_service
SyslogIdentifier=ci_service

[Install]
WantedBy=multi-user.target
