# config/demo/templates/openwrt/custom_ap/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# these services do not run on dumb APs
for i in firewall dnsmasq odhcpd; do
  if /etc/init.d/"$i" enabled; then
    /etc/init.d/"$i" disable
    /etc/init.d/"$i" stop
  fi
done

exit 0
